# ops/systemd/sandy-automation.service
[Unit]
Description=Sandy Chaos automation full-pass runner (orchestrator + autospawn + dispatch)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/projects/sandy-chaos
## Secrets loading order:
# 1) repo-local .secrets.local.env (preferred, gitignored)
# 2) legacy ~/.config/sandy-chaos/automation.env (fallback)
EnvironmentFile=-%h/projects/sandy-chaos/.secrets.local.env
EnvironmentFile=-%h/.config/sandy-chaos/automation.env
Environment=PATH=%h/.npm-global/bin:%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/python3 %h/projects/sandy-chaos/scripts/hourly_openclaw_cycle.py --dispatch-limit 1 --agent claude
StandardOutput=append:%h/.local/state/sandy-chaos/automation.log
StandardError=append:%h/.local/state/sandy-chaos/automation.log
