# ops/systemd/sandy-automation.service
[Unit]
Description=Sandy Chaos automation full-pass runner (orchestrator + autospawn + dispatch)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/dev/sandy-chaos
# Required: export OPENCLAW_TELEGRAM_BOT_TOKEN in EnvironmentFile below.
EnvironmentFile=%h/.config/sandy-chaos/automation.env
Environment=PATH=%h/.npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/python3 %h/dev/sandy-chaos/scripts/self_improve.py full-pass --scheduler host-cron --with-orchestration --with-dispatch --dispatch-limit 3 --send-telegram
StandardOutput=append:%h/.local/state/sandy-chaos/automation.log
StandardError=append:%h/.local/state/sandy-chaos/automation.log
